# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; path to your socket file

[supervisord]
# logfile=/tmp/supervisord.log ; supervisord log file
# logfile_maxbytes=50MB        ; maximum size of logfile before rotation
# logfile_backups=10           ; number of backed up logfiles
loglevel=info                ; info, debug, warn, trace
pidfile=/tmp/supervisord.pid ; pidfile location
nodaemon=false               ; run supervisord as a daemon
minfds=1024                  ; number of startup file descriptors
minprocs=200                 ; number of process descriptors
# user=yourusername            ; setuid to this UNIX account at startup
childlogdir=/tmp             ; where child log files will live

# [rpcinterface:supervisor]
# supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# [supervisorctl]
# serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

# [program:gunicorn]
# command=/opt/venv/bin/gunicorn -w 2 -b 0.0.0.0:9000 --forwarded-allow-ips=* --log-level DEBUG --timeout 120 --graceful-timeout 120 dataset_retriever.wsgi
# # directory=/path/to/your/app
# # user=yourusername
# autostart=true
# # autorestart=true
# stdout_logfile=/dev/fd/1  ; Redirect stdout to console
# stderr_logfile=/dev/fd/2  ; Redirect stderr to console
# stdout_logfile_maxbytes=0 ; Disable log rotation (optional)
# stderr_logfile_maxbytes=0 ; Disable log rotation (optional)

[program:celery]
command=/opt/venv/bin/celery -A dataset_retriever worker --loglevel=DEBUG --concurrency=1 -n worker1@%%h worker_hijack_root_logger=False worker_redirect_stdouts=True worker_redirect_stdouts_level=DEBUG
# directory=/path/to/your/app
# user=yourusername
autostart=true
# autorestart=true
stdout_logfile=/dev/fd/1  ; Redirect stdout to console
stderr_logfile=/dev/fd/2  ; Redirect stderr to console
stdout_logfile_maxbytes=0 ; Disable log rotation (optional)
stderr_logfile_maxbytes=0 ; Disable log rotation (optional)

[program:celerybeat]
command=/opt/venv/bin/celery -A dataset_retriever beat --loglevel=DEBUG --scheduler=django_celery_beat.schedulers.DatabaseScheduler worker_hijack_root_logger=False worker_redirect_stdouts=True worker_redirect_stdouts_level=DEBUG
autostart=true
# autorestart=true
stdout_logfile=/dev/fd/1  ; Redirect stdout to console
stderr_logfile=/dev/fd/2  ; Redirect stderr to console
stdout_logfile_maxbytes=0 ; Disable log rotation (optional)
stderr_logfile_maxbytes=0 ; Disable log rotation (optional)
